# node2/geth_boot.service
[Unit]
Description=Ethereum go client

[Service]
Type=simple
ExecStart=./start_nohup.sh

[Install]
WantedBy=default.target